# etc/profile.d/openssl.csh
if ( ! $?MANPATH ) setenv MANPATH ""
setenv MANPATH "${MANPATH}:/usr/ssl/man"
